# debian-tei-oxgarage/debian/postrm
#!/bin/sh
if which invoke-rc.d >/dev/null 2>&1; then
  invoke-rc.d tomcat7 restart
else
  /etc/init.d/tomcat7 restart
fi

